# ssh/com.zig.sshd-alt-port.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.zig.sshd-alt-port — second sshd listener on :2222.

  WHY: macOS's stock SSH runs via launchd socket-activation (Apple's
  /System/Library/LaunchDaemons/ssh.plist owns :22). In that mode, the
  `Port` directive in /etc/ssh/sshd_config is IGNORED — launchd passes
  pre-bound sockets to sshd. Adding additional listen ports requires
  another launchd job that owns the additional socket.

  This plist is installed at /Library/LaunchDaemons/ (user-supplied
  launch daemons; no SIP / Full Disk Access required to manage). It
  runs sshd in inetd-compatibility mode the same way Apple's does,
  just on :2222 instead of :22.

  RATIONALE: needed for zig-computer (tag:server) → pico (user-owned)
  SSH. Tailscale SSH grammar has no syntax for tag-src → user-owned-dst
  (see zig-zone runbook + bead dotfiles-wzh); tailscaled only intercepts
  :22 on the tailnet IP. Sending SSH to pico:2222 bypasses tailscaled,
  falls through to this sshd directly, authenticates via SSH key.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zig.sshd-alt-port</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/sbin/sshd</string>
        <string>-i</string>
    </array>

    <key>Sockets</key>
    <dict>
        <key>Listeners</key>
        <dict>
            <key>SockServiceName</key>
            <string>2222</string>
            <key>SockType</key>
            <string>stream</string>
        </dict>
    </dict>

    <key>inetdCompatibility</key>
    <dict>
        <key>Wait</key>
        <false/>
    </dict>

    <key>StandardErrorPath</key>
    <string>/var/log/sshd-alt-port.log</string>
</dict>
</plist>
