# ssh/com.zig.sshd-alt-port.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  com.zig.sshd-alt-port — second sshd listener on :2222.

  WHY: macOS's stock SSH runs via launchd socket-activation (Apple's
  /System/Library/LaunchDaemons/ssh.plist owns :22). In that mode, the
  `Port` directive in /etc/ssh/sshd_config is IGNORED — launchd passes
  pre-bound sockets to sshd. Adding additional listen ports requires
  either another socket-activated launchd job OR a daemon-mode sshd
  that owns its own socket. We use the latter (daemon mode) because
  socket-activation in a user-supplied plist hits hidden Apple
  entitlement requirements — runs=0 even with sockets bound (verified
  2026-05-24 on macOS Sonoma 14.4).

  This plist runs sshd in daemon mode (`-D`), letting sshd itself
  open its `-p 2222` listener. launchd just supervises (KeepAlive,
  RunAtLoad). No inetdCompatibility, no special entitlements.

  RATIONALE: needed for zig-computer (tag:server) → pico (user-owned)
  SSH. Tailscale SSH grammar has no syntax for tag-src → user-owned-dst
  (see zig-zone runbook + bead dotfiles-wzh); tailscaled only intercepts
  :22 on the tailnet IP. Sending SSH to pico:2222 bypasses tailscaled,
  falls through to this sshd directly, authenticates via SSH key.

  Host keys come from /etc/ssh/ssh_host_*_key (Apple's sshd-keygen-
  wrapper generates them on first SSH-enable; should already exist).
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zig.sshd-alt-port</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/sbin/sshd</string>
        <string>-D</string>
        <string>-p</string>
        <string>2222</string>
        <string>-e</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardErrorPath</key>
    <string>/var/log/sshd-alt-port.log</string>
    <key>StandardOutPath</key>
    <string>/var/log/sshd-alt-port.log</string>
</dict>
</plist>
