# ss14-wrapper/com.zig.ss14-wrapper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.zig.ss14-wrapper</string>
	<key>ProgramArguments</key>
	<array>
		<string>USER_HOME_PLACEHOLDER/ss14-wrapper/ss14-wrapper</string>
	</array>
	<key>WorkingDirectory</key>
	<string>USER_HOME_PLACEHOLDER/ss14-wrapper</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>SS14_WRAPPER_LISTEN</key>
		<string>TAILSCALE_IP_PLACEHOLDER:11214</string>
		<key>SS14_WRAPPER_UPSTREAM</key>
		<string>127.0.0.1:11215</string>
		<key>SS14_WRAPPER_SOCK</key>
		<string>/tmp/ss14-wrapper.sock</string>
		<key>SS14_WRAPPER_TTL_SECONDS</key>
		<string>120</string>
		<key>SS14_WRAPPER_LOG</key>
		<string>/tmp/ss14-wrapper.log</string>
	</dict>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/ss14-wrapper.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/ss14-wrapper.log</string>
</dict>
</plist>
