# a1111/com.zig.a1111.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.zig.a1111</string>
	<key>ProgramArguments</key>
	<array>
		<string>USER_HOME_PLACEHOLDER/stable-diffusion-webui/webui.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>USER_HOME_PLACEHOLDER/stable-diffusion-webui</string>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/opt/homebrew/bin:/opt/homebrew/sbin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
		<key>HOME</key>
		<string>USER_HOME_PLACEHOLDER</string>
	</dict>
	<!-- ON-DEMAND ONLY. Do NOT auto-start. SD-WebUI loads its diffusion
	     checkpoint into RAM at launch (~17 GB resident BEFORE any image
	     request), starving LLM inference (mlx_lm.server, ollama) of
	     unified memory. Start manually via `sd-up` / `launchctl kickstart
	     gui/$UID/com.zig.a1111` only when you need image generation;
	     stop via `sd-down` / `launchctl stop com.zig.a1111` when done. -->
	<key>RunAtLoad</key>
	<false/>
	<!-- KeepAlive intentionally omitted — was forcing respawn after manual
	     kill, defeating on-demand control. -->
	<key>StandardErrorPath</key>
	<string>/tmp/a1111.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/a1111.log</string>
</dict>
</plist>
